# cron.d/indexers
00 00 * * * __USER__ __BASEDIR__/bin/search --index
